# fractal.service
[Unit]
Description=Mandelbrot Fractal Generator - VULNERABLE (DoS Testing Only)
After=network.target
Wants=network.target

[Service]
Type=simple
User=kali
WorkingDirectory=/var/www/html/fractal
ExecStart=/usr/bin/python3 /var/www/html/fractal/fractal_app.py
Restart=always
RestartSec=10

# ============================================
# INTENTIONALLY VULNERABLE SETTINGS
# ============================================

# NO MEMORY LIMITS - Allows giant fractals to consume all RAM
# MemoryMax=512M  ← COMMENTED OUT = UNLIMITED

# NO CPU LIMITS - Can use 100% of all cores
# CPUQuota=50%    ← COMMENTED OUT = UNLIMITED

# UNLIMITED PROCESSES - Can spawn many workers
# TasksMax=5      ← COMMENTED OUT = UNLIMITED

# LONG TIMEOUTS - Keeps hanging processes alive
TimeoutStartSec=300      # 5 minutes to start
TimeoutStopSec=300       # 5 minutes to stop
TimeoutSec=300          # 5 minutes total

# Allow unlimited file descriptors
LimitNOFILE=65535

# Keep process alive even if misbehaving
KillMode=process        # Don't kill child processes
Restart=always

# No protection settings
ProtectSystem=no
PrivateTmp=no
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
