# fractal.service
[Unit]
Description=Mandelbrot Fractal Generator - VULNERABLE (DoS Testing Only)
After=network.target
Wants=network.target

[Service]
Type=simple
User=kali
WorkingDirectory=/var/www/html/fractal
ExecStart=/usr/local/bin/gunicorn --workers 10 --bind 127.0.0.1:5000 fractal_app:app
Restart=always
RestartSec=10

# ============================================
# INTENTIONALLY VULNERABLE SETTINGS
# ============================================

# NO MEMORY LIMITS - Allows giant fractals to consume all RAM
# MemoryMax=512M

# NO CPU LIMITS - Can use 100% of all cores
# CPUQuota=50%

# UNLIMITED PROCESSES - Can spawn many workers
# TasksMax=5

# LONG TIMEOUTS - Keeps hanging processes alive
TimeoutStartSec=300
TimeoutStopSec=300
TimeoutSec=300

# Allow unlimited file descriptors
LimitNOFILE=65535

# Keep process alive even if misbehaving
KillMode=process
Restart=always

# No protection settings
ProtectSystem=no
PrivateTmp=no
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
